# deploy/dryft-bot.service
[Unit]
Description=Dryft Telegram Bot
After=network.target

[Service]
Type=simple
User=dryft
Group=dryft
WorkingDirectory=/opt/dryft
EnvironmentFile=/opt/dryft/.env
ExecStart=/opt/dryft/venv/bin/python telegram_bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/dryft/state /opt/dryft/.cache
PrivateTmp=true

[Install]
WantedBy=multi-user.target
